# model/config/model.service
[Unit]
Description=Public Turnilo
After=syslog.target

[Service]
ExecStart=turnilo --examples
User=www-data
Group=www-data
Restart=always
StandardError=syslog
NotifyAccess=all
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target