# model/config/model.service
[Unit]
Description=Language Tool API service
After=syslog.target

[Service]
ExecStart=java -cp /etc/api-endpoint/lt/LanguageTool-5.6/languagetool-server.jar org.languagetool.server.HTTPServer --port 8081 --allow-origin
User=www-data
Group=www-data
Restart=always
StandardError=syslog
NotifyAccess=all
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target