# model/config/model.service
[Unit]
Description=Public Turnilo service
After=syslog.target

[Service]
ExecStart=turnilo run-config /etc/api-endpoint/config.yaml
User=www-data
Group=www-data
Restart=always
StandardError=syslog
NotifyAccess=all
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target